# project/cfg/BoardConfig_IPC/overlay/overlay-luckfox-ubuntu/etc/rc.local
#!/bin/bash

/usr/bin/filesystem_resize.sh
/etc/init.d/S50usbdevice start
luckfox-config load
if [ -n "$(hwclock | grep "invalid")" ]; then
	date -s 2024-01-01
	hwclock -w
fi
