# project/cfg/BoardConfig_IPC/overlay/overlay-luckfox-alpine-spi_nand/etc/init.d/ubiattach
#!/sbin/openrc-run

depend() {
    need devfs
    before localmount
}

start() {
    ebegin "Attaching UBI volumes"
    ubiattach /dev/ubi_ctrl -m 5
    ubiattach /dev/ubi_ctrl -m 4
    eend $? "Failed to attach UBI volumes"
}
