# sysdrv/tools/board/eudev/rules/61-usbdevice.rules
SUBSYSTEM=="android_usb",ACTION=="change",RUN+="/usr/bin/usbdevice %E{DEVPATH} %E{USB_STATE}"
